# mlforkids-numbers/supervisord.conf
[program:ml4k-numbers]
command=/usr/bin/python3 server.py
autostart=true
autorestart=true
stderr_logfile=/var/log/ml4k-numbers.err.log
stdout_logfile=/var/log/ml4k-numbers.out.log
